# examples/resolvloser.path
[Unit]
Description=Watch for changes in /etc/resolv.conf

[Path]
PathChanged=/etc/resolv.conf

[Install]
WantedBy=network-online.target
